# mariadb/rootfs/etc/services.d/mysql-general-log/run
#!/usr/bin/with-contenv sh

exec logutil-service -f /var/run/s6/mysql-general-log-fifo $HENCE_APP_VOL_PREFIX/logs/mysql-general-logs
