# varnish/rootfs/etc/services.d/varnish/finish
#!/usr/bin/execlineb -P
exit 0 #exiting this way, otherwise may kill containser before child process starts correctly
