# mariadb/rootfs/etc/services.d/mysql-slow-query-log/run
#!/usr/bin/with-contenv sh

exec logutil-service -f /var/run/s6/mysql-slow-query-log-fifo $HENCE_APP_VOL_PREFIX/logs/mysql-slow-query-logs
