# mariadb/rootfs/etc/services.d/mysql-error-log/run
#!/usr/bin/with-contenv sh

exec logutil-service -f /var/run/s6/mysql-error-log-fifo $HENCE_APP_VOL_PREFIX/logs/mysql-error-logs
