# memcached/rootfs/etc/services.d/memcached-error-log/run
#!/usr/bin/execlineb -P

import -u HENCE_APP_VOL_PREFIX logutil-service -f /var/run/s6/memcached-error-log-fifo ${HENCE_APP_VOL_PREFIX}/logs/memcached-error-logs
