# mailhog/rootfs/etc/services.d/mailhog/run
#!/usr/bin/with-contenv sh

set -e
source $HENCE_PREFIX/hence-utils.sh

set -- /go/bin/MailHog

chown -R $HENCE_APP_USER:$HENCE_APP_USER \
  $GOPATH || true

print_welcome_page

echo '4'
exec "$@"
